# dns-updater.timer
[Unit]
Description=Timer for dns-updater

[Timer]
OnUnitActiveSec=20m

[Install]
WantedBy=timers.target
